# tests/config2.conf
[program:program_name]
command
directory=/etc
priority=999
