# tests/config1.conf
[program:program_name]
command=/bin/ls
directory=/etc
priority=999
